# pkg/PLIST
@mandir /usr/local/man
bin/miki
man/man8/miki.8
@cwd /
etc/rc.d/miki
