# pkg/PLIST
bin/miki
man/man8/miki.8
@exec makewhatis /usr/local/man
@unexec makewhatis /usr/local/man
@cwd /
etc/rc.d/miki
